# debian/temboard.cron.d
3-53/10 * * * * temboard /usr/bin/temboard-metric-aggregator >/dev/null 2>&1
